# cron/crontab
# Finanzen Crontab - Alle automatisierten Tasks
# Logs werden in /app/data/logs/ geschrieben

# FinTS-Daten täglich um 6:00 Uhr abrufen
0 6 * * * cd /app && python3 scripts/fetch_fints.py >> /app/data/logs/fetch_fints.log 2>&1

# Transaktionen täglich um 7:00 Uhr kategorisieren
0 7 * * * cd /app && python3 scripts/categorize.py >> /app/data/logs/categorize.log 2>&1

# PDFs alle 2 Stunden verarbeiten
0 */2 * * * cd /app && python3 scripts/parse_pdfs.py >> /app/data/logs/parse_pdfs.log 2>&1

# Postbank-CSV (data/inbox/**/*.csv) alle 2 Stunden
5 */2 * * * cd /app && python3 scripts/import_postbank_csv.py >> /app/data/logs/import_postbank_csv.log 2>&1

# Monitoring täglich um 8:00 Uhr
0 8 * * * /bin/bash /app/cron/monitor.sh >> /app/data/logs/monitor.log 2>&1

# Alte Logs aufräumen (älter als 90 Tage) - wöchentlich sonntags um 3:00
0 3 * * 0 find /app/data/logs -name "*.log" -mtime +90 -delete

# Leere Zeile am Ende für korrekte Crontab-Formatierung
